# config/supervisor/dj_init.conf
[program:dj_init]
command=/home/dj_init/dj_init.com/dj_init/scripts/gunicorn.bash
user=dj_init
